# archiso-profile/airootfs/home/live/.config/autostart/ahmados-session-defaults.desktop
[Desktop Entry]
Type=Application
Name=Lumina-OS Session Defaults
Exec=/usr/local/bin/lumina-apply-session-defaults
OnlyShowIn=KDE;
NoDisplay=true
X-KDE-autostart-phase=2
