# archiso-profile/airootfs/home/live/.config/autostart/ahmados-welcome.desktop
[Desktop Entry]
Type=Application
Name=AhmadOS Welcome
Exec=/usr/local/bin/ahmados-welcome --once
OnlyShowIn=KDE;
NoDisplay=true
X-KDE-autostart-phase=3
